# danted-dongle0.service
[Unit]
Description=SOCKS5 proxy - Dongle 0 (True) port 1080
After=network.target dongle-init.service

[Service]
Type=simple
ExecStartPre=/bin/sh -c 'fuser -k 1080/tcp 2>/dev/null; sleep 2; true'
ExecStart=/usr/sbin/danted -N 1 -f /etc/danted-dongle0.conf
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
